# infra/systemd/ae-wiki-worker.service
[Unit]
Description=ae-wiki-agent worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ae-wiki-agent
EnvironmentFile=/opt/ae-wiki-agent/.env
ExecStart=/usr/local/bin/bun src/cli.ts jobs:worker
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
